# systemd/pm-clob-collector.service
[Unit]
Description=Polymarket CLOB WebSocket Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pm-research
Group=pm-research
WorkingDirectory=/opt/pm-research
EnvironmentFile=/etc/pm-research/.env
ExecStart=/opt/pm-research/venv/bin/python -m pm_research.entrypoints.pm_clob_collector
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pm-clob-collector

[Install]
WantedBy=multi-user.target
